# config/runit/nmbd/run
#!/bin/sh -x
sleep 6
exec nmbd --foreground -d 2